# etc/service/hbase-master/run
#!/bin/bash
/usr/lib/hbase/bin/hbase master start
