# etc/service/hbase-thrift/run
#!/bin/bash
/usr/lib/hbase/bin/hbase thrift start
